# bind/centos6/cfg_files/supervisord.d/cron.ini
[program:crond]
command=/usr/sbin/crond -n
