# lighttpd/centos7/cfg_files/supervisord.d/sshd.ini
[program:sshd]
command=/etc/init.d/sshd -D -E /var/log/sshd.log
