# lighttpd/centos8/cfg_files/supervisord.d/lighttpd.ini
[program:lighttpd]
command=/usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D
