# lighttpd/centos8/cfg_files/supervisord.d/rsyslog.ini
[program:rsyslog]
command=/sbin/rsyslogd -i /var/run/syslogd.pid -c 5 -n
