# bind/centos7/cfg_files/supervisord.d/named.ini
[program:named]
command=/etc/init.d/named start
